# systemd/dirt.service.d/override.conf
# Drop-in for /etc/systemd/system/dirt.service.d/override.conf
#
# The dirt-camera daemon runs as a user service (systemctl --user), so its
# socket lives under $XDG_RUNTIME_DIR — which is /run/user/<uid> for the
# owning user. System-level services do NOT inherit XDG_RUNTIME_DIR, so
# dirt's Python capture client would otherwise fall back to /tmp/dirt-camera.sock
# and fail to reach the daemon. Point it explicitly.
#
# Install:
#   sudo mkdir -p /etc/systemd/system/dirt.service.d
#   sudo cp systemd/dirt.service.d/override.conf /etc/systemd/system/dirt.service.d/
#   sudo systemctl daemon-reload && sudo systemctl restart dirt

[Service]
Environment=DIRT_CAMERA_SOCKET=/run/user/1000/dirt-camera.sock
