# systemd/dirt.service
[Unit]
Description=Dirt Home Grow Monitor
After=network.target

[Service]
Type=simple
User=akcom
WorkingDirectory=/home/akcom/code/dirt
ExecStart=/home/akcom/code/dirt/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
