# systemd/dirt-camera-agent.service
[Unit]
Description=Dirt camera agent - edge snapshot capture and cloud upload
Documentation=file://%h/code/dirt/docs/commands.md
After=default.target network-online.target dirt-camera.service
Wants=network-online.target
# The camera-agent publishes scoped camera assets only; it is not a full gateway.
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=dirt-service-alert@%n.service

[Service]
Type=simple
WorkingDirectory=%h/code/dirt
EnvironmentFile=-%h/code/dirt/.env
EnvironmentFile=%h/code/dirt/.env.dirt2-camera-agent
ExecStart=%h/.local/bin/uv run --package dirt-camera-agent python -m dirt_camera_agent.main
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10s
SuccessExitStatus=143

NoNewPrivileges=true

[Install]
WantedBy=default.target
