# systemd/dirt-hwd.service
[Unit]
Description=Dirt hardware daemon — humidifier, archive, device watchdog, ingest API
Documentation=file://%h/code/dirt/CLAUDE.md
After=default.target network-online.target
Wants=network-online.target
# Ordered after the camera daemon — capture loop talks to its socket.
After=dirt-camera.service
# Outer crash-loop budget: 5 nonzero exits in 5 minutes → systemd leaves
# the unit `failed`. Reset with `systemctl --user reset-failed dirt-hwd`.
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=dirt-service-alert@%n.service

[Service]
Type=simple
WorkingDirectory=%h/code/dirt
EnvironmentFile=%h/code/dirt/.env
# `uv run --package dirt-hwd` activates the workspace venv with dirt-hwd's
# deps on PATH. python -m dirt_hwd.main invokes the uvicorn entrypoint,
# which binds :8000.
ExecStart=%h/.local/bin/uv run --package dirt-hwd python -m dirt_hwd.main
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
# Allow async context managers (capture_loop daemon RPC, kasa plug
# teardown, HTTP client pools) time to unwind before SIGKILL.
TimeoutStopSec=10s
# uvicorn translates SIGTERM into a graceful shutdown that exits with
# status 143 (128 + SIGTERM) — from systemd's POV that's `code=exited` not
# `code=killed`, so the default SIGTERM-is-success rule doesn't apply and
# OnFailure= would page on every planned restart. Whitelist 143 explicitly.
SuccessExitStatus=143

# Minimal hardening — don't break USB/camera socket access.
NoNewPrivileges=true

[Install]
WantedBy=default.target
