# systemd/dirt-daily-report.service
[Unit]
Description=Dirt daily report run (photos + sensor synthesis + Telegram)
Documentation=file://%h/code/dirt/AGENTS.md
After=dirt-camera.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/code/dirt
EnvironmentFile=%h/code/dirt/.env
ExecStart=%h/.local/bin/uv run --package dirt-shared python %h/code/dirt/scripts/daily_report
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900s
